# .arch/usr/lib/systemd/system/volta-plex.automount
[Unit]
Description=Automount network shares on boot

[Automount]
Where=/mnt/volta/plex

[Install]
WantedBy=multi-user.target

