# .arch/usr/lib/systemd/system/volta-backups.automount
[Unit]
Description=Automount network shares on boot

[Automount]
Where=/mnt/volta/backups

[Install]
WantedBy=multi-user.target

